# rootfs/etc/cron.d/remove_tomcat_logs
# Remove log files older than 15 days.
0 0 * * * root (find /opt/tomcat/logs -name \*log\* -mtime +15 -exec rm -f {} \; >/dev/null 2>&1 )
